# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/bin/php-fpm --nodaemonize

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
